# modules/pje/files/regional/utilitarios/jboss/pje-jboss.cron
#!/bin/bash

# Arquivo: pje-jboss.cron
# Objetivo: apagar arquivos .upload e .crl do diretório /tmp de cada VM JBoss, mantendo apenas os ultimos 7 dias
# Localizacao: /etc/cron.daily

## Remocao de arquivos .upload
find /tmp -name *.upload -mtime +7 -exec rm {} \;

## Remocao de arquivos .crl
find /tmp -name *.crl -mtime +7 -exec rm {} \;

## Remocao de arquivos .tmp
find /tmp -name *.tmp -mtime +7 -exec rm {} \;